# wayland-mouse.service
[Unit]
Description=wayland-mouse — Mac-like pointer + wheel acceleration (Wayland/libinput)
Documentation=https://github.com/monfa-red/wayland-mouse
After=multi-user.target
StartLimitIntervalSec=0

[Service]
Type=simple
ExecStartPre=-/sbin/modprobe uinput
ExecStart=/usr/local/bin/wayland-mouse run
Restart=always
RestartSec=2
Nice=-10

[Install]
WantedBy=multi-user.target
